# scripts/hand-tracking.service
[Unit]
Description=Hand Tracking Service V3 (DepthAI)
After=network.target multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=nvidia
WorkingDirectory=/home/nvidia/dev/HandTrackingV3/build
ExecStartPre=/home/nvidia/dev/HandTrackingV3/scripts/jetson_max_performance.sh
ExecStart=/home/nvidia/dev/HandTrackingV3/build/HandTrackingService
Restart=always
RestartSec=5

# Environment
Environment="DISPLAY=:0"
Environment="LD_LIBRARY_PATH=/usr/local/lib:/usr/lib/aarch64-linux-gnu"

# Performance Tuning - Real-time priority
Nice=-20
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=90

# Resource Limits
LimitRTPRIO=95
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
